# lesson_4/part_1/bird_watching.timer
[Unit]
Description=Activate bird_watching.service every minute

[Timer]
OnUnitInactiveSec=1m

[Install]
WantedBy=timers.target
